# salt/states/phassa/moniek.service
[Unit]
Description=moniek manages the financial accounts for giedo

[Service]
ExecStart=/home/sys-moniek/bin/moniek
EnvironmentFile=-/etc/default/moniek
User=sys-moniek
Type=notify

[Install]
WantedBy=multi-user.target

