# salt/states/sankhara/cilia-tunnel.service
[Unit]
Description=Links Cilia socket on phassa and sankhara

[Service]
ExecStart=/root/bin/cilia-tunnel
