# salt/states/sankhara/moniek-tunnel.service
[Unit]
Description=Links Moniek socket on phassa and sankhara

[Service]
ExecStart=/root/bin/moniek-tunnel
